# etc/local.rc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
Usage:
  - Place this file in ~/local/etc/local.rc.plist
  - Create a symbolic link using **full absolute paths**:
        ln -s /Users/yourname/local/etc/local.rc.plist /Users/yourname/Library/LaunchAgents/local.rc.plist

        Use absolute paths instead of relative ones.
        Relative symlinks may cause `launchctl` to fail with:
        "Unload failed: 5: Input/output error"
  - Load with:
        launchctl load ~/Library/LaunchAgents/local.rc.plist
  - Unload with:
        launchctl unload ~/Library/LaunchAgents/local.rc.plist
  - Reload with:
        launchctl kickstart -k gui/$(id -u)/local.rc.plist
  - Your init script should be at:
        ~/local/etc/rc (executable)
  - Log output:
        stdout: ~/local/var/log/local.rc.out
        stderr: ~/local/var/log/local.rc.err
-->
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.rc</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/yourname/local/etc/rc</string>
  </array>

  <!-- Run script immediately after user login -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Trigger script after system wakes from sleep (macOS Ventura and later) -->
  <key>LaunchEvents</key>
  <dict>
    <key>com.apple.powermanagement.systempowerstate</key>
    <dict>
      <key>Conditions</key>
      <dict>
        <key>Wake</key>
        <true/>
      </dict>
    </dict>
  </dict>

  <!-- Optional: log stdout and stderr to temporary files -->
  <key>StandardOutPath</key>
  <string>/Users/yourname/local/var/log/local.rc.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/yourname/local/var/log/local.rc.err</string>

  <!-- Do not restart the job if it exits -->
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
